# services/agetty/finish
#!/bin/sh
exec utmpset -w -- "${PWD##*-}"
