# services/agetty/run
#!/bin/sh
port="${PWD##*-}"
[ "$port" = tty1 ] && options=--noclear
exec chpst -P -- agetty $options -- "$port" 38400 linux
